# server-go/deploy/git-ai.service
[Unit]
Description=Git-AI Go Server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=git-ai
Group=git-ai

WorkingDirectory=/opt/git-ai/server-go/current
EnvironmentFile=/opt/git-ai/.env
ExecStart=/opt/git-ai/server-go/current/git-ai-server
ExecStartPre=/usr/bin/test -x /opt/git-ai/server-go/current/git-ai-server
ExecStartPre=/usr/bin/test -f /opt/git-ai/.env

Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
KillSignal=SIGTERM

LimitNOFILE=65536
UMask=0077
SyslogIdentifier=git-ai

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/opt/git-ai/logs

[Install]
WantedBy=multi-user.target
